# supervisord.conf
[supervisord]
nodaemon=true

[program:ts3-afk-bot]
command=python /app/main.py
autostart=true
autorestart=true
stderr_logfile=/var/log/ts3-afk-bot.err.log
stdout_logfile=/var/log/ts3-afk-bot.out.log
